# steamos-extension-loader/overlayfs/usr/lib/systemd/system/rauc.service.d/steam-extension-loader.conf
[Unit]
Conflicts=systemd-sysext.service
After=systemd-sysext.service
