# steamos-extension-clean-games/overlayfs/usr/lib/systemd/user/steamos-extension-clean-games.service
[Unit]
Description=Clean Steam Installation

[Service]
Type=oneshot
ExecStart=/usr/bin/steamos-extension-clean-games
