# steamos-extension-loader/overlayfs/usr/lib/systemd/system/steamos-extension-loader-system-cleaner.timer
[Unit]
Description=Cleanup System Extension Loader

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
