# steamos-extension-loader/overlayfs/usr/lib/systemd/system/steamos-extension-loader-system-cleaner.service
[Unit]
Description=Cleanup System Extension Loader

[Service]
Type=oneshot
ExecStart=/usr/sbin/steamos-extension-loader-system-cleaner
