# src/steamos-extension-performance-tuning/overlayfs/usr/lib/udev/rules.d/99-steamos-extension-performance-tuning-zram.rules
ACTION=="add|change", SUBSYSTEM=="block", KERNEL=="zram[0-9]|zram[0-9][0-9]", RUN+="/usr/sbin/steamos-extension-performance-tuning-dispatch zram %k"
