# src/steamos-extension-performance-tuning/overlayfs/usr/lib/udev/rules.d/99-steamos-extension-performance-tuning-mmcblk.rules
ACTION=="add|change", SUBSYSTEM=="block", KERNEL=="mmcblk[0-9]|mmcblk[0-9][0-9]", RUN+="/usr/sbin/steamos-extension-performance-tuning-dispatch mmcblk %k"
