# src/steamos-extension-update-flatpak/overlayfs/usr/lib/systemd/user/steamos-extension-update-user-flatpak.service
[Unit]
Description=Update User Flatpak Packages

[Service]
Type=oneshot
ExecStart=/usr/bin/steamos-extension-update-user-flatpak
