# steamos-extension-disable-mitigations/overlayfs/usr/lib/systemd/system/steamos-extension-kernel-cmdline.service.d/steamos-extension-disable-mitigations.conf
[Unit]
Wants=steamos-extension-disable-mitigations.service
After=steamos-extension-disable-mitigations.service
