# src/steamos-extension-loader/overlayfs/usr/lib/systemd/system/systemd-sysext.service.d/steamos-extension-loader.conf
[Unit]
Conflicts=rauc.service
Before=rauc.service
