# src/steamos-extension-loader/overlayfs/usr/lib/systemd/user/steamos-extension-loader-user-cleaner.timer
[Unit]
Description=Cleanup Extension Loader

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
