# steamos-extension-update-flatpak/overlayfs/usr/lib/systemd/user/steamos-extension-update-user-flatpak.timer
[Unit]
Description=Update User Flatpak Packages

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
