# steamos-extension-loader/overlayfs/usr/lib/systemd/user/steamos-extension-loader-user-cleaner.service
[Unit]
Description=Cleanup User Extension Loader

[Service]
Type=oneshot
ExecStart=/usr/bin/steamos-extension-loader-user-cleaner
