# steamos-extension-performance-tuning/overlayfs/usr/lib/systemd/system/steamos-extension-kernel-cmdline.service.d/steamos-extension-performance-tuning.conf
[Unit]
Wants=steamos-extension-performance-tuning.service
After=steamos-extension-performance-tuning.service
