# steamos-extension-retain-boot/overlayfs/usr/lib/systemd/system/steamos-extension-kernel-cmdline.service.d/steamos-extension-retain-boot.conf
[Unit]
Wants=steamos-extension-retain-boot.service
After=steamos-extension-retain-boot.service
