# src/steamos-extension-clean-games/overlayfs/usr/lib/systemd/user/steamos-extension-clean-games.timer
[Unit]
Description=Clean Steam Installation

[Timer]
OnCalendar=hourly
Persistent=false

[Install]
WantedBy=timers.target
