# steamos-extension-compat-tools/overlayfs/usr/lib/systemd/user/steamos-extension-compat-tools-update-luxtorpeda.service
[Unit]
Description=Update Luxtorpeda Compatibility Tool

[Service]
Type=oneshot
ExecStart=/usr/bin/steamos-extension-compat-tools-update-luxtorpeda
