# steamos-extension-compat-tools/overlayfs/usr/lib/systemd/user/steamos-extension-compat-tools-update-boxtron.service
[Unit]
Description=Update Boxtron Compatibility Tool

[Service]
Type=oneshot
ExecStart=/usr/bin/steamos-extension-compat-tools-update-boxtron
