# steamos-extension-update-flatpak/overlayfs/usr/lib/systemd/system/steamos-extension-update-flatpak.service
[Unit]
Description=Update Flatpak Applications
Requires=network.target network-online.target systemd-sysext.service

[Service]
Type=simple
ExecStart=/usr/sbin/steamos-extension-update-flatpak
