# steamos-extension-performance-tuning/overlayfs/usr/lib/udev/rules.d/99-steamos-extension-performance-tuning-sd.rules
ACTION=="add|change", SUBSYSTEM=="block", KERNEL=="sd[a-zA-Z]|sd[a-zA-Z][a-zA-z]", RUN+="/usr/sbin/steamos-extension-performance-tuning-dispatch sd %k"
