# steamos-extension-performance-tuning/overlayfs/usr/lib/udev/rules.d/99-steamos-extension-performance-tuning-wifi.rules
ACTION=="add|change", SUBSYSTEM=="net", KERNEL=="wlp*|wlan*", RUN+="/usr/sbin/steamos-extension-performance-tuning-dispatch wifi %k"
