# src/steamos-extension-compat-tools/overlayfs/usr/lib/systemd/user/steamos-extension-compat-tools-update-boxtron.timer
[Unit]
Description=Update Boxtron Compatibility Tool

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
